# deploy/templates/hubot-restart.cron.j2
MAILTO=""
0 * * * * root systemctl restart hubot.service >> /home/{{ ansible_user_id }}/crons.log